# src/leiningen/new/pedestal_micro/supervisord.conf
[supervisord]
nodaemon=true

[program:{{name}}]
command=/bin/bash -c "java -jar /srv/{{name}}.jar"

